# systemd/systemdunits/newsbeuter-reload.service
[Unit]
Description=reloads newsbueter 

[Service]
ExecStart=/usr/bin/newsbueter -x reload 

[Install]
WantedBy=default.target
